# linux/etc/init.d/shadowsocks
#!/sbin/runscript
# Copyright 1999-2014 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# 
# Author: sloger
# Description: Shadowsocks client script file for OpenRC
#

profile=${RC_SVCNAME#*.}

config_file=${config_file:-/opt/shadowsocks/config.json}
run_dir=${run_dir:-/run/shadowsocks}

command="/opt/shadowsocks/shadowsocks-Linux64bit"
command_args="-c ${config_file}"
command_background="true"

pidfile=${run_dir}/${profile}.pid

depend() {
    use net logger
}

start_pre() {
    checkpath -d -m 0750 "${run_dir}"
    if [ ! -f ${config_file} ]; then
        eerror "Missing configuration file ${config_file}"
        return 1
    fi
}
